# ground-station/systemd/dump978-fa.service
[Unit]
Description=dump978-fa UAT decoder (FlyCatcher_UAT)
After=network.target

[Service]
Type=simple
User=blastly
# Before first reboot: device=1 by index. After reboot: serial=UAT978 takes effect.
# Post-reboot the service will find it by serial automatically.
ExecStart=/usr/local/bin/dump978-fa \
    --sdr driver=rtlsdr,serial=UAT978 \
    --sdr-auto-gain \
    --json-port 30978
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
