# ground-station/systemd/vhf-cleanup.service
[Unit]
Description=VHF audio cleanup (prune oldest when >1TB)

[Service]
Type=oneshot
User=blastly
ExecStart=/home/blastly/scripts/vhf-cleanup.sh
